# deployment/nodetypes/dev/etc/cron.daily/svn-add-node-logs
#!/bin/bash

svn add --parents --auto-props /node/log/*
svn ci -m "[dev-node] Add daily node logs" /node/log
