# deployment/root.d/eci-secca/etc/profile.d/nodeenv.sh

PATH=$PATH:/node/type/bin:/node/type/sbin:/node/local/bin

if [ -f /node/.noderc ]; then
    . /node/.noderc
fi
