# deployment/root.d/eci-secca/etc/profile.d/su-keephist.sh
function su() {
    if [ $# = 0 ]; then
        sudo su --login
    else
        /bin/su --login "$@"
    fi
}
