# deployment/node/sdu/sdu-dat-seccaj/seccaj.conf
description "seccaj (dev) application domain"
author "Lenik <xjl@99jsj.com>"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel[!2345]

emits seccaj

console output
task

exec seccaj start

# Not work...
# pre-stop exec seccaj stop
